# crontab/lse-docker.crontab
# =============================================================================
# Crontab для хоста с Docker-контейнером lse-bot (пример: user ai8049520).
#
# Установка (заменит ВЕСЬ crontab пользователя — сначала сделайте бэкап):
#   crontab -l > ~/crontab.backup.$(date +%F).txt
#   crontab /home/ai8049520/lse/crontab/lse-docker.crontab
#
# Другой пользователь / каталог:
#   sed 's|/home/ai8049520/lse|/home/YOU/lse|g' lse/crontab/lse-docker.crontab | crontab -
#
# Премаркет: 30 12,15 MSK — окно до 09:30 ET (EDT). Старая строка 15 17 MSK летом = уже RTH, без Telegram.
# =============================================================================

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 10,12,14,16,18,20,22 * * 1-5 docker exec lse-bot python scripts/update_prices_cron.py >> /dev/null 2>&1
0 22 * * 0,6 docker exec lse-bot python scripts/update_prices_cron.py >> /dev/null 2>&1
10 22 * * * docker exec lse-bot python scripts/update_rsi_local.py >> /dev/null 2>&1
0 19 * * 1-5 docker exec lse-bot python update_finviz_data.py >> /dev/null 2>&1
0 16,19,22 * * 1-5 docker exec lse-bot python scripts/trading_cycle_cron.py >> /home/ai8049520/lse/logs/cron_trading_cycle.log 2>&1
*/5 * * * 1-5 docker exec lse-bot python scripts/send_sndk_signal_cron.py >> /home/ai8049520/lse/logs/cron_sndk_signal.log 2>&1
30 12,15 * * 1-5 docker exec lse-bot python scripts/premarket_cron.py >> /home/ai8049520/lse/logs/premarket_cron.log 2>&1
# Earnings autoprep (полный цикл documents): KB yfinance → sync → ingest → LLM extract → readiness JSON.
# ML train/spillover/shadow: run_earnings_ml_refresh ниже. Analyzer только прогресс готовности.
15 */2 * * * flock -n /tmp/lse_earnings_autoprep.lock docker exec lse-bot python scripts/run_earnings_intelligence_autoprep.py >> /home/ai8049520/lse/logs/earnings_autoprep.log 2>&1
# Phase C: Telegram brief after new LLM extraction (past earnings)
0 9,18 * * 1-6 flock -n /tmp/lse_earnings_brief_notify.lock docker exec lse-bot python scripts/notify_earnings_post_event_brief.py >> /home/ai8049520/lse/logs/earnings_post_event_brief.log 2>&1
# Trust Arbiter: morning digest 09:05 MSK (06:05 UTC) — LSE Trust operator screen
5 6 * * 1-6 flock -n /tmp/lse_trust_digest.lock docker exec lse-bot python scripts/notify_trust_digest.py --refresh >> /home/ai8049520/lse/logs/trust_digest.log 2>&1
# Earnings post-mortem + trust metrics (after quotes/outcomes; complements shadow in ml refresh)
10 0 * * 2-6 flock -n /tmp/lse_earnings_postmortem.lock docker exec lse-bot python scripts/run_earnings_postmortem_cron.py >> /home/ai8049520/lse/logs/earnings_postmortem.log 2>&1
# Unified ML refresh dispatcher (data-driven train trigger + ml_contours_status.json). See docs/ML_UNIFIED_RETRAIN_FRAMEWORK.md
15 */6 * * * flock -n /tmp/lse_ml_refresh_dispatcher.lock docker exec lse-bot python scripts/run_ml_refresh_dispatcher.py >> /home/ai8049520/lse/logs/ml_refresh_dispatcher.log 2>&1
# Earnings ML (legacy slot — optional; dispatcher covers */6 poll). Full train: 23:52 below.
# 30 */6 * * * flock -n /tmp/lse_earnings_ml_refresh.lock docker exec lse-bot python scripts/run_earnings_ml_refresh.py >> ...
# News ingest: free stack (RSS ЦБ + yfinance earnings). NewsAPI отключён (426 на free tier).
*/15 * * * * flock -n /tmp/lse_news_core_fast.lock docker exec lse-bot python scripts/fetch_news_cron.py --mode core-fast >> /home/ai8049520/lse/logs/news_fetch.log 2>&1
# Ticker news: Yahoo yfinance.get_news по GAME_5M / portfolio (без ключа).
5 */2 * * * flock -n /tmp/lse_news_tickers.lock docker exec lse-bot python scripts/fetch_news_cron.py --mode tickers >> /home/ai8049520/lse/logs/news_fetch.log 2>&1
# Investing: economic calendar JSON (работает) + HTML news (может 403 без прокси).
0 */2 * * * flock -n /tmp/lse_news_investing.lock docker exec lse-bot python scripts/fetch_news_cron.py --mode investing >> /home/ai8049520/lse/logs/news_fetch.log 2>&1
# NewsAPI — выключено (NEWSAPI 426 Upgrade Required). Включить при платном плане:
# 10 */2 * * * flock -n /tmp/lse_news_newsapi.lock docker exec lse-bot python scripts/fetch_news_cron.py --mode newsapi >> /home/ai8049520/lse/logs/news_fetch.log 2>&1
10 * * * * docker exec lse-bot python scripts/sync_vector_kb_cron.py >> /home/ai8049520/lse/logs/sync_vector_kb.log 2>&1
20 * * * * docker exec lse-bot python scripts/add_sentiment_to_news_cron.py >> /home/ai8049520/lse/logs/add_sentiment_to_news.log 2>&1
0 4 * * * docker exec lse-bot python scripts/analyze_event_outcomes_cron.py >> /home/ai8049520/lse/logs/analyze_event_outcomes.log 2>&1
30 4 * * * docker exec lse-bot python scripts/cleanup_calendar_noise.py --execute >> /home/ai8049520/lse/logs/cleanup_calendar_noise.log 2>&1
45 * * * * docker exec lse-bot python scripts/cron_watchdog.py --execute >> /home/ai8049520/lse/logs/cron_watchdog.log 2>&1
20 15 * * 1-5 flock -n /tmp/lse_gap_forecast_premarket.lock docker exec lse-bot python scripts/ingest_game5m_gap_forecast.py --ensure-table --phase premarket >> /home/ai8049520/lse/logs/game5m_gap_forecast.log 2>&1
40 16 * * 1-5 flock -n /tmp/lse_gap_forecast_open.lock docker exec lse-bot python scripts/ingest_game5m_gap_forecast.py --ensure-table --phase open >> /home/ai8049520/lse/logs/game5m_gap_forecast.log 2>&1
30 23 * * 1-5 docker exec -e ANALYZER_AUTOTUNE_APPLY=0 lse-bot python scripts/analyzer_autotune.py --days 5 --url http://127.0.0.1:8080/api/analyzer >> /home/ai8049520/lse/logs/autotune.log 2>&1
25 23 * * * flock -n /tmp/lse_market_bars_intraday.lock docker exec lse-bot python scripts/ingest_market_bars_intraday.py >> /home/ai8049520/lse/logs/cron_market_bars_intraday.log 2>&1
# market_regime_daily (SPY/QQQ/DIA/VIX from quotes) — before event-reaction build
32 23 * * 1-5 flock -n /tmp/lse_market_regime.lock docker exec lse-bot python scripts/ingest_market_regime_daily.py >> /home/ai8049520/lse/logs/event_regime.log 2>&1
# Event / earnings reaction dataset: скелет из KB → разметка из daily quotes (до GAME_5M ML).
# Product advisory model uses expanded history: v0_expanded_baseline + quotes_regime_v1.
33 23 * * 1-5 flock /tmp/lse_erd_nightly.lock docker exec lse-bot python scripts/build_event_reaction_dataset.py --from-kb-earnings --dataset-version v0_expanded_baseline --include-earnings-universe >> /home/ai8049520/lse/logs/event_reaction_build.log 2>&1
35 23 * * 1-5 flock /tmp/lse_erd_nightly.lock docker exec lse-bot python scripts/backfill_event_reaction_labeling.py --dataset-version v0_expanded_baseline --only-outcomes --force-outcomes --include-earnings-universe --limit 800 >> /home/ai8049520/lse/logs/event_reaction_outcomes.log 2>&1
39 23 * * 1-5 flock /tmp/lse_erd_nightly.lock docker exec -e EVENT_REACTION_FEATURE_BUILDER_VERSION=quotes_regime_v1 lse-bot python scripts/backfill_event_reaction_labeling.py --dataset-version v0_expanded_baseline --include-all-symbols --limit 8000 >> /home/ai8049520/lse/logs/event_reaction_labeling.log 2>&1
# Earnings intelligence layer: overwrite features on earnings universe (superset of quotes_regime_v1).
43 23 * * 1-5 flock /tmp/lse_erd_nightly.lock docker exec -e EVENT_REACTION_FEATURE_BUILDER_VERSION=quotes_regime_earnings_v1 lse-bot python scripts/backfill_event_reaction_labeling.py --dataset-version v0_expanded_baseline --only-features --force-features --include-earnings-universe --limit 800 >> /home/ai8049520/lse/logs/event_reaction_earnings_labeling.log 2>&1
45 23 * * 1-5 flock -n /tmp/lse_erd_labeling_gaps.lock docker exec lse-bot python scripts/check_erd_labeling_gaps_cron.py --execute --telegram >> /home/ai8049520/lse/logs/erd_labeling_gaps.log 2>&1
# GAME_5M ML datasets: stuck + continuation (B2); bar v2 apply off (список A).
48 23 * * 1-5 flock -n /tmp/lse_game5m_daily_ml.lock docker exec -e DAILY_ML_RUN_ENTRY_BAR_V2_APPLY=0 lse-bot python scripts/run_daily_game5m_ml_pipeline.py >> /home/ai8049520/lse/logs/game5m_daily_ml_pipeline.log 2>&1
# Open-path labels (before nightly ML refresh).
52 23 * * 1-5 flock -n /tmp/lse_open_path_labels.lock docker exec lse-bot python scripts/label_open_path_scenarios.py --ensure-table >> /home/ai8049520/lse/logs/open_path_labels.log 2>&1
# Unified ML nightly refresh: game5m_entry, open_path, event_reaction, earnings_grid, portfolio (+ 20d in portfolio refresh). See docs/ML_AND_DECISION_ARCHITECTURE.md §5.3
55 23 * * 1-5 flock -n /tmp/lse_ml_refresh_nightly.lock docker exec -e ML_READINESS_TRAIN_MODE=full lse-bot python scripts/run_ml_refresh_dispatcher.py --slot nightly >> /home/ai8049520/lse/logs/ml_refresh_nightly.log 2>&1
# ML train readiness: skip A contours (game5m v1, event_reaction); B earnings grid stays on; portfolio 20d nightly write via PORTFOLIO_CATBOOST_20D_NIGHTLY_WRITE.
58 23 * * 1-5 flock -n /tmp/lse_ml_train_readiness.lock docker exec -e ML_READINESS_SKIP_EVENT_REACTION=1 -e ML_READINESS_SKIP_EARNINGS_INTELLIGENCE=0 -e ML_READINESS_SKIP_GAME5M=1 -e EVENT_REACTION_DATASET_VERSION=v0_expanded_baseline -e EVENT_REACTION_FEATURE_BUILDER_VERSION=quotes_regime_earnings_v1 lse-bot python scripts/run_ml_train_readiness_cron.py >> /home/ai8049520/lse/logs/ml_train_readiness_cron.log 2>&1
# Portfolio 20d session probe (log_only observe vs game) after readiness.
5 0 * * 2-6 flock -n /tmp/lse_portfolio_20d_probe.lock docker exec lse-bot python scripts/run_portfolio_20d_session_probe.py >> /home/ai8049520/lse/logs/portfolio_20d_session_probe.log 2>&1
# Monthly portfolio allocation report (~6m path + capture + 20d) → /portfolio/daily above charts. 1-е 07:20 MSK.
20 7 1 * * flock -n /tmp/lse_portfolio_prospect_monthly.lock docker exec lse-bot python scripts/run_portfolio_prospect_monthly_review.py >> /home/ai8049520/lse/logs/portfolio_prospect_monthly.log 2>&1
# Legacy per-contour nightly (replaced by dispatcher --slot nightly 23:47):
# 46 open_path | 51 event_reaction --full | 52 earnings_grid --full
# Open-path */6 poll — via run_ml_refresh_dispatcher.py (15 */6)
# Open-path weekly full train — FROZEN (список A).
# 48 23 * * 0 flock -n /tmp/lse_open_path_ml_full.lock docker exec -e ML_READINESS_OPEN_PATH_TRAIN_MODE=full -e OPEN_PATH_ML_REFRESH_APPLY_DATA=1 -e OPEN_PATH_ML_CONTINUOUS_TRAIN=1 lse-bot python scripts/run_open_path_ml_refresh.py --full >> /home/ai8049520/lse/logs/open_path_ml_full.log 2>&1
# Weekly ML full refresh: multiday_lr, recovery, gap_forecast (+ open_path, game5m_entry, game5m_entry_bar_v2).
5 6 * * 0 flock -n /tmp/lse_ml_refresh_weekly.lock docker exec -e ML_READINESS_TRAIN_MODE=full lse-bot python scripts/run_ml_refresh_dispatcher.py --slot weekly_full >> /home/ai8049520/lse/logs/ml_refresh_weekly.log 2>&1
# Weekly earnings prod eval (B5 — labels/grid research).
0 6 * * 0 flock -n /tmp/lse_earnings_prod_eval.lock docker exec -e ML_READINESS_TRAIN_MODE=full -e EARNINGS_ML_REFRESH_APPLY_DATA=1 -e EARNINGS_ML_REFRESH_INCREMENTAL_TRAIN=1 lse-bot python scripts/run_earnings_intelligence_prod_eval.py >> /home/ai8049520/lse/logs/earnings_prod_eval.log 2>&1
# Quotes gaps for ERD (no_quotes backlog) — воскресенье: deep backfill earnings universe.
25 5 * * 0 flock -n /tmp/lse_erd_seed_quotes_weekly.lock docker exec lse-bot python scripts/seed_quotes_for_event_reaction_dataset.py --dataset-version v0_expanded_baseline --include-earnings-universe --all-symbols --days 450 >> /home/ai8049520/lse/logs/erd_seed_quotes.log 2>&1
# Отчёт качества данных без LLM и без профилирования local/datasets/*.csv.
53 23 * * 1-5 flock -n /tmp/lse_ml_data_quality_report.lock docker exec lse-bot python scripts/run_ml_data_quality_report.py --no-default-datasets --json-out /app/logs/ml/ml_data_quality/report_daily.json >> /home/ai8049520/lse/logs/ml_data_quality_report_cron.log 2>&1
# Recovery D4a stats — FROZEN (список A).
# 54 23 * * 1-5 flock -n /tmp/lse_recovery_d4a_stats.lock docker exec lse-bot python scripts/run_recovery_d4a_stats_cron.py >> /home/ai8049520/lse/logs/recovery_d4a_stats_cron.log 2>&1
# Decision stack mirror telemetry (phase 3): weekly shadow diff, no RESOLVE gate.
10 6 * * 0 flock -n /tmp/lse_decision_stack_mirror.lock docker exec lse-bot python scripts/report_decision_stack_mirror.py --days 14 --strategy ALL --json-out /app/logs/ml/ml_data_quality/last_decision_stack_mirror_report.json >> /home/ai8049520/lse/logs/decision_stack_mirror.log 2>&1
# Options sentiment gate shadow (phase 4): closed GAME_5M vs would_downgrade + live scan; log_only, no apply.
15 6 * * 0 flock -n /tmp/lse_options_gate_shadow.lock docker exec lse-bot python scripts/analyze_options_gate_shadow.py --days 28 >> /home/ai8049520/lse/logs/options_gate_shadow.log 2>&1
# Multiday walk-forward OOS (GAME_5M tickers) — weekly before RTH.
25 6 * * 0 flock -n /tmp/lse_multiday_wf.lock docker exec lse-bot python3 /app/scripts/run_multiday_wf_game5m.py >> /home/ai8049520/lse/logs/multiday_wf_game5m.log 2>&1
# Post-RTH GAME_5M session review + gap rolling metrics (after US close ~23:30 MSK).
35 23 * * 1-5 flock -n /tmp/lse_game5m_daily_review.lock /home/ai8049520/lse/scripts/cron_game5m_daily_session_review.sh >> /home/ai8049520/lse/logs/game5m_daily_review.log 2>&1
# yfinance option OI daily snapshot → options_chain_oi_snapshot (Money Map history slider). docs/OPTIONS_MONEY_MAP.md
30 23 * * 1-5 flock -n /tmp/lse_options_oi_snap.lock /home/ai8049520/lse/scripts/cron_options_chain_oi.sh
# Money Map PCR stats from cron OI snapshots (p25/p75 thresholds vs wireframe). docs/OPTIONS_MONEY_MAP.md § Cron stats
45 23 * * 1-5 flock -n /tmp/lse_options_map_stats.lock /home/ai8049520/lse/scripts/cron_options_map_stats.sh
# GAME_5M replay proposals (weekly, exit-family only) + observe active experiment (weekdays after RTH).
40 23 * * 1-5 flock -n /tmp/lse_game5m_tuning_observe.lock /home/ai8049520/lse/scripts/cron_game5m_tuning.sh observe >> /home/ai8049520/lse/logs/game5m_tuning.log 2>&1
35 6 * * 0 flock -n /tmp/lse_game5m_tuning_propose.lock /home/ai8049520/lse/scripts/cron_game5m_tuning.sh propose >> /home/ai8049520/lse/logs/game5m_tuning.log 2>&1
40 6 * * 0 flock -n /tmp/lse_weekly_game5m_tactic_review.lock /home/ai8049520/lse/scripts/cron_weekly_game5m_tactic_review.sh >> /home/ai8049520/lse/logs/weekly_game5m_tactic_review.log 2>&1
# GAME_5M light analyzer 7d (B6 observability). docs/GAME_5M_SESSION_ANALYSIS_PLAYBOOK.md §5
45 6 * * 0 flock -n /tmp/lse_game5m_light_analyzer_weekly.lock /home/ai8049520/lse/scripts/cron_weekly_game5m_light_analyzer.sh >> /home/ai8049520/lse/logs/game5m_light_analyzer_weekly.log 2>&1
# B2 continuation go/no-go (список B): 14.07 primary, 21.07 backup — docs/GAME_5M_B_LIST_RUN_SCHEDULE.md
15 7 14 7 * flock -n /tmp/lse_b2_gonogo.lock /home/ai8049520/lse/scripts/cron_game5m_b2_continuation_gonogo.sh
15 7 21 7 * flock -n /tmp/lse_b2_gonogo.lock /home/ai8049520/lse/scripts/cron_game5m_b2_continuation_gonogo.sh
# Gap forecast rolling metrics refresh (after open ingest; updates last_gap_forecast_metrics.json).
50 16 * * 1-5 flock -n /tmp/lse_gap_analyze_daily.lock docker exec lse-bot python3 /app/scripts/analyze_game5m_gap_forecast.py --days 90 --json-metrics-out /app/logs/ml/ml_data_quality/last_gap_forecast_metrics.json >> /home/ai8049520/lse/logs/game5m_gap_forecast.log 2>&1

# Premarket ML feature snapshots: script writes only during PRE_MARKET.
10 12,14 * * 1-5 flock -n /tmp/lse_premarket_features.lock docker exec lse-bot python scripts/ingest_premarket_daily_features.py --ensure-table >> /home/ai8049520/lse/logs/premarket_features.log 2>&1
30 15 * * 1-5 flock -n /tmp/lse_premarket_features.lock docker exec lse-bot python scripts/ingest_premarket_daily_features.py --ensure-table >> /home/ai8049520/lse/logs/premarket_features.log 2>&1
20 16 * * 1-5 flock -n /tmp/lse_premarket_features.lock docker exec lse-bot python scripts/ingest_premarket_daily_features.py --ensure-table >> /home/ai8049520/lse/logs/premarket_features.log 2>&1
